# updates/configs/plexrequestsnet.service
[Unit]
Description=Ombi - PMS Requests System
After=network-online.target

[Service]
User=root
Group=root
WorkingDirectory=/opt/plexrequest.net/
ExecStart=/opt/plexrequest.net/Ombi
Type=simple
TimeoutStopSec=30
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
